# services/tidal-watchdog.service
[Unit]
Description=Tidal Connect Watchdog services. Stops and starts tidal when device is connected/disconnected
Wants=tidal-watchdog.timer

[Service]
Type=oneshot
User=root

ExecStart=/usr/bin/tidal-watchdog.sh

[Install]
WantedBy=multi-user.target






