# services/tidal-watchdog.service
[Unit]
Description=Tidal Connect Watchdog services. Stops and starts tidal when amp goes off and on
Wants=tidal-watchdog.timer
ConditionPathExists=/etc/tidal/config.json

[Service]
Type=oneshot
User=root

ExecStart=/usr/bin/tidal-watchdog.sh

[Install]
WantedBy=multi-user.target






